# example/alpine-nginx-ruby/services.d/rackup/run
#!/usr/bin/execlineb -P

define custom_runner "/srv/code/mobingi-runapp"
with-contenv
cd /srv/code

ifthenelse
{ test -f $custom_runner }
{
  exec $custom_runner
}
{
  exec rackup -o 0.0.0.0
}
