# example/alpine-nginx-ruby/services.d/rackup/run
#!/usr/bin/execlineb -P

define custom_runner "/srv/code/mobingi-runapp.sh"
cd /srv/code
with-contenv

ifthenelse
{ test -f $custom_runner }
{ exec $custom_runner }
{ exec ./bundle_bin/rackup -o 0.0.0.0 }


