# example/nginx-statichtml/services.d/nginx/run
#!/usr/bin/execlineb -P

# set specific user and group
s6-envuidgid nginx
nginx -g "daemon off;"
