# example/alpine-nginx-ruby/services.d/puma/run
#!/usr/bin/execlineb -P

cd /srv/code
with-contenv
exec ./bundle_bin/puma
