# example/ubuntu-apache2-php7-with-cron/services.d/apache2/run
#!/usr/bin/with-contenv bash

source /etc/apache2/envvars
/usr/sbin/apache2 -DFOREGROUND
