# example/alpine-nginx-ruby/services.d/redis/run
#!/usr/bin/execlineb -P

# set specific user and group
s6-envuidgid redis
redis-server /etc/redis.conf --daemonize no --logfile /dev/stdout
